# tools/cron/crontab
23      3,21    *       *       *       letsencrypt > /dev/null
45      23      2       *       *       /usr/bin/php    BASE_PATH/cron/update-ip2location-db.php
